# postgres.conf
[program:postgres]
user=root
command=/init-pg.sh
autorestart=true
stopsignal=INT
